# install/example.cron
*/5 * * * * wget -q http://your.url.here/api.php?type=action&id=mailqueue > /dev/null 2>&1